# src/redis-server.conf
[program:redis-server]
command = redis-server /opt/docker-registry/redis.conf
